# tools/fosbuilder/src/_rootfs/etc/init.d/aaaa-sysinit-depmod
#!/sbin/openrc-run

# Test depmod ; Alex BOTZUNG

description="Run depmod."

depend() {
	before devfs
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Probing modules (with depmod)"
	depmod
	eend $?
}
