# tools/fosbuilder/src/_rootfs/etc/init.d/aaaa-default-keylayout
#!/sbin/openrc-run

# Set the keyboard layout ; Alex BOTZUNG

description="Set the correct keyboard layout"

depend() {
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Setting keyboard layout"
	bash /bin/FOS-changekbd
	eend $?
}
