# tools/fosbuilder/src/_rootfs/etc/init.d/aaab-default-network
#!/sbin/openrc-run

# Start FOG ; Alex BOTZUNG

description="Start the network"

depend() {
	before aaaa-sysinit-depmod
    provide network
    provide net
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Start the network"
	bash /bin/FOS-network
	eend $?
}
