# tools/fosbuilder/src/_rootfs/etc/init.d/aaac-default-vncserver
#!/sbin/openrc-run

# Spawn early shell ; Alex BOTZUNG

description="Start VNC server"

depend() {
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Start VNC server"
	bash /bin/FOS-vncserver &
	eend $?
}
