# tools/fosbuilder/src/_rootfs/etc/init.d/aaad-default-changerootpwd
#!/sbin/openrc-run

# Change root password ; Alex BOTZUNG

description="Change root password"

depend() {
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Change root password"
	bash /bin/FOS-ChangeRootPWD
	eend $?
}
