# tools/fosbuilder/src/_rootfs/etc/init.d/aaae-default-rebootin24hrs
#!/sbin/openrc-run

# Change root password ; Alex BOTZUNG

description="Reboot FOS in 24 hrs"

depend() {
	true
}

start() {
	ebegin "Reboot FOS in 24 hrs"
	busybox reboot -f -d 86400 &
	eend $?
}
