# tools/fosbuilder/src/_rootfs/etc/init.d/aaae-default-earlyshell
#!/sbin/openrc-run

# Spawn early shell ; Alex BOTZUNG

description="Spawn an early shell"

depend() {
	keyword -prefix -lxc -docker
}

start() {
	ebegin "Spawn early shell"
	bash /bin/FOS-earlyShell
	eend $?
}
