# tools/fosbuilder/src/_rootfs/etc/init.d/aaaf-default-FOS-InstallAPK
#!/sbin/openrc-run

# Change root password ; Alex BOTZUNG

description="Install optional apk packages"

depend() {
    after network
}

start() {
	ebegin "Install optional apk packages"
	/bin/FOS-InstallAPK --run-as-service
	eend $?
}
